# packages/NetTraRec/NetTraRec-0.1.7.tar.gz/NetTraRec-0.1.7/NetRecorder/nettrarec.raw
#!/bin/sh
# chkconfig: 35 20 80
# description: netrec-server
~/miniconda/bin/netrec-server
